# recipes-core/amixeron/files/amixeron.service
[Unit]
Description=Turn On Amixer IEC958 

[Service]
ExecStart=/usr/sbin/amixeron.sh

[Install]
WantedBy=multi-user.target
